# luci-app-4eburnet/files/etc/hotplug.d/iface/40-4eburnet
#!/bin/sh
# 4eburNet — восстановление TPROXY routing при реконнекте WAN
# Запускается procd при каждом ifup/ifdown событии

# Только при поднятии интерфейса
[ "$ACTION" = "ifup" ] || exit 0

# Только WAN интерфейсы (PPPoE, DHCP, статика, LTE)
case "$INTERFACE" in
    wan|wan6|pppoe-wan|pppoe*|wan_*|*wan*|lte*|4g*)
        ;;
    *)
        exit 0
        ;;
esac

# 4eburNet запущен?
[ -f /var/run/4eburnet.pid ] || exit 0
PID=$(cat /var/run/4eburnet.pid 2>/dev/null)
[ -n "$PID" ] || exit 0
kill -0 "$PID" 2>/dev/null || exit 0
# Защита от PID reuse: проверить имя процесса
[ "$(cat /proc/$PID/comm 2>/dev/null)" = "4eburnetd" ] || exit 0

logger -t 4eburnet "WAN $INTERFACE поднят, проверяем routing..."

# Параметры из UCI (fallback на дефолты из constants.h)
FWMARK=$(uci -q get 4eburnet.main.fwmark 2>/dev/null || echo "0x01")
TABLE=$(uci -q get 4eburnet.main.route_table 2>/dev/null || echo "100")
PRIO=$(uci -q get 4eburnet.main.route_prio 2>/dev/null || echo "1000")

# Восстановить ip rule fwmark → table (TPROXY)
# grep -qE: ядро может показать 0x1 или 0x00000001
ip rule show 2>/dev/null | grep -qE "fwmark 0x(0{0,7}1) " || {
    ip rule add fwmark "$FWMARK" table "$TABLE" prio "$PRIO" 2>/dev/null && \
        logger -t 4eburnet "ip rule восстановлен (fwmark $FWMARK → table $TABLE prio $PRIO)"
}

# Восстановить ip route local 0.0.0.0/0 dev lo table
ip route show table "$TABLE" 2>/dev/null | grep -q "local" || {
    ip route add local 0.0.0.0/0 dev lo table "$TABLE" 2>/dev/null && \
        logger -t 4eburnet "ip route восстановлен (local → table $TABLE)"
}

# IPv6: ip -6 rule fwmark → table
ip -6 rule show 2>/dev/null | grep -qE "fwmark 0x(0{0,7}1) " || {
    ip -6 rule add fwmark "$FWMARK" table "$TABLE" prio "$PRIO" 2>/dev/null
}

# IPv6: ip -6 route local ::/0 dev lo table
ip -6 route show table "$TABLE" 2>/dev/null | grep -q "local" || {
    ip -6 route add local ::/0 dev lo table "$TABLE" 2>/dev/null
}

# Перечитать конфиг (новый WAN IP → обновить провайдеров)
/usr/sbin/4eburnetd --ipc reload 2>/dev/null && \
    logger -t 4eburnet "конфиг перечитан после WAN $INTERFACE"

exit 0
